# oomd.service
[Unit]
Description=Facebook OOMD Service starter

[Service]
User=root
ExecStart=/usr/local/bin/oomd_bin

[Install]
WantedBy=multi-user.target
